# _etc/udev/rules.d/99-hdmi-plug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/jnsbl/.Xauthority", RUN+="/home/jnsbl/.dotfiles/bin/update-lockscreen"
